# ansible/crons/refresh_burp_servers.crontab
MS_PREFIX=/srv/makina-states
D=/etc/default/makina-states
0 */1 *  * * root sh -c "if [ -f ${D} ];then . ${D};fi;${MS_PREFIX}/ansible/crons/refresh_burp_servers.sh"
# vim:set et sts=4 ts=4 tw=80 ft=crontab:

