# files/etc/systemd/system/lxc-setup.service
# MANAGED VIA SALT -- DO NOT EDIT
[Unit]
Description=lxc-setup
After=sysinit.target
DefaultDependencies=no
Before=network-online.target network.target multi-user.target dbus.service
After=systemd-tmpfiles-setup.service

[Service]
Type=oneshot
ExecStart=/usr/bin/lxc-setup.sh systemd
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
