# files/etc/systemd/system/lxc-net-makina.service
# MANAGED VIA SALT -- DO NOT EDIT
#
[Unit]
Description=lxc-net-makina

[Service]
Type=oneshot
ExecStart=/usr/bin/lxc-net-makina.sh start
ExecStop=/usr/bin/lxc-net-makina.sh stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
