# files/etc/systemd/system/lxc-setup.service
# MANAGED VIA SALT -- DO NOT EDIT
[Unit]
Description=lxc-setup

[Service]
Type=oneshot
ExecStart=/usr/bin/ms-lxc-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
