# files/etc/systemd/system/lxc-stop.service
# MANAGED VIA SALT -- DO NOT EDIT
[Unit]
Description=lxc-stop
Before=shutdown.target reboot.target halt.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/bin/true
ExecStop=/usr/bin/lxc-stop.sh

[Install]
WantedBy=shutdown.target
