# files/etc/systemd/system/ms_supervisor.service
# MANAGED VIA SALT -- DO NOT EDIT
[Unit]
Description=ms_supervisor
After=syslog.target network.target

[Service]
ExecStart=/etc/init.d/ms_supervisor start --nodaemon
ExecStop=/etc/init.d/ms_supervisor stop
ExecReload=/etc/init.d/ms_supervisor reload
[Install]
WantedBy=multi-user.target
