# files/etc/cron.d/slapdle
0 2 * * * root /etc/slapd_le.sh &>/dev/null
