# files/etc/systemd/system/lxc.service.d/lxc.conf
[Unit]
After=syslog.target network.target lxc-net.service lxc-net-makina.service
Wants=lxc-net.service lxc-net-makina.service
