# ansible/roles/ms_services_base_ntpd/templates/etc/cron.d/ntpsync
0 2 * * * root if [ -x /sbin/ntp-sync.sh ];then /sbin/ntp-sync.sh;fi
