# files/etc/init/circusd.conf
# MANAGED VIA SALT, DO NOT EDIT
{% set data = salt['mc_utils.json_load'](data) %}
console output
start on filesystem and net-device-up IFACE=lo
stop on shutdown

script
exec su -c "/usr/bin/circus.sh start"
end script

pre-stop script 
/usr/bin/circus.sh stop
end script

post-stop script 
/usr/bin/circus.sh stop
end script
