# files/etc/init/circusd.conf
# MANAGED VIA SALT, DO NOT EDIT
{% set data = salt['mc_utils.json_load'](data) %}
console output
start on filesystem and net-device-up IFACE=lo
stop on shutdown
script
. /etc/profile || /bin/true
. {{ data.venv }}/bin/activate
exec {{ data.venv }}/bin/circusd --log-output {{data.log}} --pidfile {{data.pidf}} {{data.conf}}
end script
#post-stop script
#. /etc/profile || /bin/true
#. "{{ data.venv }}/bin/activate"
#ps aux |grep circus|grep python|awk '{print $2}'|xargs kill -9 || /bin/true
#end script
