# files/etc/init/circusd.conf
# MANAGED VIA SALT, DO NOT EDIT
console output
start on filesystem and net-device-up IFACE=lo
stop on shutdown
script
. {{ venv }}/bin/activate
exec {{ venv }}/bin/circusd --log-output {{log}} --pidfile {{pidf}} {{ locs['conf_dir'] }}/circus/circusd.ini
end script
