# files/etc/init/circusd.conf
# MANAGED VIA SALT, DO NOT EDIT
{% set data = salt['mc_utils.json_load'](data) %}
console output
start on filesystem and net-device-up IFACE=lo
stop on shutdown
script
. /etc/profile || /bin/true
. {{ data.venv }}/bin/activate
exec su -c "{{ data.venv }}/bin/circusd {{data.conf}}"
end script
pre-stop script
. /etc/profile || /bin/true
. "{{ data.venv }}/bin/activate"

is_lxc() {
    echo  "$(cat -e /proc/1/environ |grep container=lxc|wc -l|sed -e "s/ //g")"
}

filter_host_pids() {
    if [ "x$(is_lxc)" != "x0" ];then
        echo "${@}"
    else
        for pid in ${@};do
            if [ "x$(grep -q lxc /proc/${pid}/cgroup 2>/dev/null;echo "${?}")" != "x0" ];then
                 echo ${pid}
             fi
         done
    fi
}
{% if data.circusd.httpd %}
circusctl stop circushttpd --waiting || /bin/true
{% endif %}
{% if data.circusd.statsd %}
circusctl stop circusd-stats --waiting || /bin/true
{% endif %}
pids=$(filter_host_pids $(ps aux|grep "from circus import stats"|awk '{print $2}'))
if [ "x${pids}" != "x" ];then
  kill -9 "${pids}"
fi
pids=$(filter_host_pids $(ps aux|grep "from circusweb import"|awk '{print $2}'))
if [ "x${pids}" != "x" ];then
  kill -9 "${pids}"
fi
end script
