# files/etc/init/ms_supervisor.conf
# MANAGED VIA SALT, DO NOT EDIT
{% set data = salt['mc_utils.json_load'](data) %}
console output
start on filesystem and net-device-up IFACE=lo
stop on shutdown
script
. /etc/profile || /bin/true
export TMPDIR=/tmp
. {{ data.venv }}/bin/activate
exec su -c '{{ data.venv }}/bin/supervisord -c "{{data.conf}}" --nodaemon'
end script
