# files/etc/init/salt-master.conf
{% set salt_data = salt['mc_salt.settings']().data_mappings[service_name][salt_mode] %}
description "{{salt_data['daemon_name']}}"
{% if salt_data['init_debug'] %}
console output
{% else %}
# console output
{% endif %}
start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]
# docker containers wont support ulimit's
# https://groups.google.com/d/msg/docker-user/T45Kc9vD804/xz7UZk_cilsJ
{% if salt_data['has_filelimit'] %}
limit nofile 100000 100000
{% else %}
# limit nofile 100000 100000
{% endif %}
script
# BADIDEA restart memcache to be consistent
# if [ -e /etc/init.d/memcached ] || [ -e /etc/init/memcached.conf ];then
#   service memcached restart || /bin/true
# fi
# remove stale sockets
rm -f {{salt_data['sock_dir']}}/{master_event,publish_pull,workers}*ipc
exec {{salt_data['bin_dir']}}/{{salt_data['daemon_name']}} {% if salt_data['init_debug'] %} -lall {% endif %}
end script
post-stop script
# remove stale sockets
rm -f {{salt_data['sock_dir']}}/minion_event*ipc
end script
