# files/etc/init/uwsgi.conf
# MANAGED VIA SALT -- DO NOT EDIT
{% set data = salt['mc_uwsgi.settings']() %}
description "uWSGI Emperor"
start on runlevel [2345]
stop on runlevel [06]

# docker containers wont support ulimit's
# https://groups.google.com/d/msg/docker-user/T45Kc9vD804/xz7UZk_cilsJ
limit nofile 100000 100000
exec "/usr/bin/uwsgi.sh"
