# files/etc/profile.d/salt.sh
export PATH=/srv/salt/bin/:$PATH
export PATH={{ saltbinpath }}:$PATH
